# audio-receiver/audio-receiver.service
[Unit]
Description=Audio Stream Receiver
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/audio-receiver
ExecStart=/usr/bin/python3 /opt/audio-receiver/receiver.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
