# web-ui/web-ui.service
[Unit]
Description=Audio Archive Web UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/web-ui
ExecStart=/usr/bin/python3 /opt/web-ui/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
